# lintian-brush/tests/systemd-service-alias-without-extension/override/in/debian/foo.service
[Unit]
Description=Test Service
Alias=bar

[Service]
Type=oneshot
