# lintian-brush/tests/systemd-service-file-refers-to-obsolete-bindto/override/in/debian/foo.service
[Unit]
Description=Test Service
BindTo=bar.service

[Service]
Type=oneshot
